# packages/snipsskills/snipsskills-0.1.4.847-py2.py3-none-any.whl/snipsskills/config/systemd/snipsskills@USER.service
[Unit]
Description=Snips Skills
After=network.target

[Service]
Type=simple
User=%i
ExecStart=snips &> /dev/null & cd {{SNIPS_HOME_PATH}}; {{SNIPSSKILLS_PATH}} run

[Install]
WantedBy=multi-user.target