# packages/portinus/portinus-1.0.0.tar.gz/portinus-1.0.0/portinus/templates/monitor.timer
[Unit]
Description=Monitor timer for {{ name }} service

[Timer]
OnActiveSec=5m

[Install]
WantedBy=timers.target
