# packages/portinus/portinus-1.0.1.tar.gz/portinus-1.0.1/portinus/templates/monitor.service
[Unit]
Description=Monitor module for {{ name }} service
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/portinus-monitor check --name {{ name }}

