# packages/snipsskills/snipsskills-0.1.4.856-py2.py3-none-any.whl/snipsskills/config/systemd/snipsskills@USER.service
[Unit]
Description=Snips Skills
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory={{SNIPS_HOME_PATH}}
ExecStart={{SNIPSSKILLS_PATH}} run

[Install]
WantedBy=multi-user.target